# examples/connect-tcp/client/qrexec-connect-tcp@2222.socket.d/git.conf
[Socket]
ListenStream=127.0.1.1:2222
FileDescriptorName=git
