# examples/connect-tcp/client/qrexec-connect-tcp@2222.socket.d/sshfs.conf
[Socket]
ListenStream=127.0.1.2:2222
FileDescriptorName=sshfs
