# examples/split-gpg-extra/client/qrexec-connect-unix@gpg-agent-extra.socket.d/@default.conf
[Socket]
ListenStream=%t/qrexec-connect-unix/S.gpg-agent.extra
FileDescriptorName=@default
