# examples/split-ssh/client/qrexec-connect-unix@gpg-agent-ssh.socket.d/@default.conf
[Socket]
ListenStream=%t/qrexec-connect-unix/S.gpg-agent.ssh
FileDescriptorName=@default
