# examples/split-ssh/server/gpg-agent-ssh.socket.d/qubes.ConnectUnix.conf
[Socket]
Symlinks=%t/qubes.ConnectUnix/gpg-agent-ssh
