# systemd-user/qrexec-connect-unix@.service
[Unit]
Description=Connect through qubes.ConnectUnix

[Service]
ExecStart=/opt/bin/qrexec-systemd-socket-activate --rpc-arg='%i' qubes.ConnectUnix
