# systemd-user/qrexec-connect-tcp@.socket
[Unit]
Description="Connect through qubes.ConnectTCP"

[Install]
WantedBy=sockets.target
