# systemd-user/qrexec-connect.service
[Unit]
Description=Bind client-service qube network addresses

[Service]
ExecStart=/opt/bin/qrexec-connect
