# systemd-user/qrexec-connect-unix@.socket
[Unit]
Description="Connect through qubes.ConnectUnix"

[Install]
WantedBy=sockets.target
