# systemd-user/qrexec-connect-tcp-bind@.socket
[Unit]
Description="Connect through qubes.BindConnectTCP"

[Install]
WantedBy=sockets.target
