# systemd-user/qrexec-connect-.socket.d/defaults.conf
[Unit]
Description=Bind client-service qube network addresses
RequiredBy=qrexec-connect.service
Before=qrexec-connect.service

[Install]
WantedBy=sockets.target
