# systemd-user/qrexec-connect-.socket.d/service.conf
[Socket]
Service=qrexec-connect.service
